# files/system/etc/systemd/system/run-media-david-home.automount
[Unit]
Description=SMB home share on NAS

[Automount]
Where=/run/media/david/home

[Install]
WantedBy=multi-user.target
